# .config/systemd/user/newsboat.service
# newsboat.service

[Unit]
Description= Newsboat autoupdate

[Service]
Type=forking
ExecStart=/usr/bin/newsboat -x reload 
